# onlinescreensaver/bin/onlinescreensaver.conf
start on started lab126_gui
stop on stopping lab126_gui

pre-start script
    source /etc/upstart/functions

    if [ -e /mnt/us/extensions/onlinescreensaver/bin/scheduler.sh ]; then
        /bin/sh /mnt/us/extensions/onlinescreensaver/bin/scheduler.sh &
    fi

    true
end script

post-stop script
    PID=$(ps xa | grep "scheduler.sh" | grep -v grep | awk '{ print $1 }')
    kill $PID || true
end script
