# onlinescreensaver/bin/onlinescreensaver.conf
start on started onlinescreensaver-mount
stop on stopping onlinescreensaver-mount

# Use 'respawn' to ensure the scheduler stays running
respawn

# The scheduler is the main process
exec /bin/sh /mnt/base-us/extensions/onlinescreensaver/bin/scheduler.sh

post-stop script
    pkill -f "scheduler.sh" || true
end script
